# autostart/compton.desktop
[Desktop Entry]
Name=compton
Type=Application
Exec=compton --vsync opengl